# systemd/huebot.service
[Unit]
Description=Hue Bot Brain
After=network.target

[Service]
ExecStart=/root/huebot/huebot.py

[Install]
WantedBy=multi-user.target
